# etc/init/amphora-agent.conf
description "Start up the Octavia Amphora Agent"

start on startup

respawn
respawn limit 2 2

exec amphora-agent --config-file /etc/octavia/amphora-agent.conf
